# udev/20-connect-wifi.rules
# Fritz! WLAN stick
ACTION=="add", SUBSYSTEMS=="usb", ATTRS{idVendor}=="057c", ATTRS{idProduct}=="8503", RUN+="/usr/bin/netctl start fvz"
