# udev/20-connect-wifi.rules
# Fritz! WLAN stick
ACTION=="add", SUBSYSTEMS=="net", ATTRS{address}=="dc:15:c8:88:1c:a0", RUN+="/usr/bin/netctl start fvz"
